# rootfs/services/nginx/run
#!/bin/sh

exec 2>&1

mkdir -p /run/nginx
exec nginx -g "daemon off;"

